# launchd/com.a2333.horizon-trigger.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.a2333.horizon-trigger</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/opt/python@3.11/libexec/bin/python3</string>
    <string>/Users/a2333/IDE/n8n/scripts/horizon_trigger_watcher.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/a2333/IDE/n8n</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>N8N_PROJECT_DIR</key>
    <string>/Users/a2333/IDE/n8n</string>
    <key>HORIZON_DIR</key>
    <string>/Users/a2333/IDE/n8n/horizon</string>
    <key>HORIZON_TRIGGER_HOST</key>
    <string>0.0.0.0</string>
    <key>HORIZON_TRIGGER_PORT</key>
    <string>8765</string>
    <key>HORIZON_CONTAINER_NAME</key>
    <string>horizon</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/a2333/IDE/n8n/logs/horizon-trigger-launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/a2333/IDE/n8n/logs/horizon-trigger-launchd.err.log</string>
</dict>
</plist>
